# .config/sh/profile.d/apps.sh
#!/bin/sh
export MAIL=geary
export BROWSER=firefox
export TERM=alacritty

# TODO: update to use terminal emacsclient
export VISUAL=emacsclient
export EDITOR=/usr/bin/vim
