# .config/sh/profile.d/apps.sh
#!/bin/sh
# [[file:../../../Bash.org::*Apps][Apps:1]]
export MAIL=geary
export BROWSER=firefox
export TERM=alacritty

# TODO: update to use terminal emacsclient
export VISUAL=emacsclient
export EDITOR=/usr/bin/vim
# Apps:1 ends here
