# .config/sh/profile.d/emacs.sh
#!/bin/sh
# [[file:../../../Bash.org::*Emacs][Emacs:1]]
export ORG_DIRECTORY=/data/org
#export ORG_AGENDA_ROOT=
#export ORG_ROAM_ROOT
# Emacs:1 ends here
