# .config/sh/profile.d/emacs.sh
#!/bin/sh
export ORG_DIRECTORY=/data/org
#export ORG_AGENDA_ROOT=
#export ORG_ROAM_ROOT
