# .config/sh/profile.d/gtk.sh
#!/bin/sh
# [[file:../../../Bash.org::*GTK][GTK:1]]
# GTK
export GTK2_RC_FILES="$HOME/.gtkrc-2.0"

#export GDK_BACKEND=wayland             # this can prevent programs from starting (e.g. chromium and electron apps). therefore, this should be set per app instead of globally.
# GTK:1 ends here
