# .config/sh/profile.d/julia.sh
#!/bin/sh
# [[file:../../../Bash.org::*Julia][Julia:1]]

# Julia:1 ends here
