# .config/sh/profile.d/xdg.sh
#!/bin/sh
# [[file:../../../Bash.org::*XDG][XDG:1]]
export XDG_SESSION_TYPE=wayland
export XDG_SESSION_DESKTOP=sway
export XDG_CURRENT_DESKTOP=sway
# XDG:1 ends here
